# mysql-supervisord.conf
;
; config for mysql
;

[program:mysql]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe
user=mysql
priority=0


;
; config for mysql
;


